# skel/etc/profile.d/brew.sh
# Homebrew (group-shared) in /home/brew/.linuxbrew
if [ -x /home/brew/.linuxbrew/bin/brew ]; then
  eval "$(/home/brew/.linuxbrew/bin/brew shellenv)"
fi
